# bin/cryptolock.service
[Unit]
Description=ethereum cryptolock service
Requires=network.target
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cryptolock
Restart=always
KillMode=process
ExecStart=/usr/bin/node /opt/cryptolock-js/bin/cryptolock.js
WorkingDirectory=/opt/cryptolock-js/bin
