# bin/cryptolock.service
[Unit]
Description=ethereum cryptolock service
Requires=network.target
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cryptolock
Restart=always
ExecStart=node bin/cryptolock.js
