# bin/slock.service
[Unit]
Description=ethereum slock service
Requires=network.target
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=slock
Restart=always
KillMode=process
ExecStart=/usr/bin/node /opt/slock-js/bin/slock.js
WorkingDirectory=/opt/slock-js/bin
